# inotify_watch.conf
[program:inotify_watch]
command=/bin/bash /tmp/watch.sh
stopwaitsecs=5
restart=always
user=root
